# PLS/CA4/Alternative/vagrant_ubuntu/salt/host1/files/spring-boot.service
[Unit]
Description=Spring Boot Application
After=network.target
Wants=network.target

[Service]
Type=simple
User=devuser
Group=developers
WorkingDirectory=/opt/apps/spring-app
Environment="JAVA_HOME=/usr/lib/jvm/java-17-openjdk-amd64"
Environment="PATH=/usr/lib/jvm/java-17-openjdk-amd64/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/apps/spring-app/gradlew bootRun
Restart=always
RestartSec=30
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spring-boot

[Install]
WantedBy=multi-user.target